# FichiersDeConfig/etc/cron.d/backup-script
0 3 * * * root ./backup.sh >> /var/log/backup_full.log 2>&1