# FichiersDeConfig/etc/cron.d/clamav-scan
0 3 * * * root clamscan --log=/var/log/clamav/scan.log /var/www /etc